# supervisord.conf
[supervisord]
nodaemon=true

[program:cron]
command=/usr/sbin/cron -f
stdout_logfile=/var/log/cron_stdout.log
stderr_logfile=/var/log/cron_stderr.log